# supervisord.conf
[supervisord]
nodaemon=true

[program: Beat]
command=celery -A tasks worker --loglevel=info -B
autorestart=true
